# supervisord-grafana.conf
[program:grafana]
directory=/usr/share/grafana
command=/usr/sbin/grafana-server
numprocs=1
autostart=true
autorestart=true
