# supervisord-telegraf.conf
[program:telegraf]
directory=/root
command=/opt/telegraf/telegraf -config telegraf.conf
numprocs=1
autostart=true
autorestart=true
