# supervisord-influx.conf
[program:influxdb]
directory=/root
command=influxd -config influx.conf.toml
numprocs=1
autostart=true
autorestart=true

